# backend/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:migrate]
command=alembic upgrade head
autorestart=false
startsecs=0
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:web]
command=bash -c "sleep 5 && uvicorn app.main:app --host 0.0.0.0 --port %(ENV_PORT)s --workers 1"
autorestart=true
priority=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:celery-worker]
command=celery -A app.tasks.celery_app worker --loglevel=info --concurrency=2
autorestart=true
startretries=3
priority=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:celery-beat]
command=celery -A app.tasks.celery_app beat --loglevel=info
autorestart=true
startretries=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
